# systemd/service-controller.service
[Unit]
Description=Service Controller from YAML
After=network.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/service-controller.py

[Install]
WantedBy=multi-user.target
